# fusion/rootfs/etc/s6-overlay/s6-rc.d/fusion/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: fusion
# Runs fusion
# ==============================================================================

bashio::log.info 'Starting fusion...'

# Check if fusion binary exists
if [ -f /usr/local/bin/fusion ]; then
    bashio::log.debug 'Fusion binary found!'
else
    bashio::log.error 'Fusion binary not found at /usr/local/bin/fusion'
    exit 1
fi

# Check if the 'password' config option is set and not empty
if bashio::config.has_value 'password'; then
    bashio::log.debug 'Password config is not empty. Setting PASSWORD environment variable...'
    export PASSWORD="$(bashio::config 'password')"
else
    bashio::log.debug 'No password set in configuration.'
fi

# Run fusion
bashio::log.debug 'Running Fusion...'
exec /usr/local/bin/fusion